# dockerfile_build/llama_openclaw.conf
[program:llama-server]
directory=/root/llama-b9660
command=/root/llama-b9660/llama-server --model /root/bge-small-zh-v1.5-f16.gguf --port 11436 --embedding --threads 1 --host 0.0.0.0 --ubatch-size 1024
stdout_logfile=/var/log/llama.log
stderr_logfile=/var/log/llama.log
autostart=true
autorestart=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

[program:openclaw-gateway]
command=openclaw gateway
stdout_logfile=/var/log/openclaw.log
stderr_logfile=/var/log/openclaw.log
autostart=true
autorestart=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
